# setup/aws-efs.service
[Unit]
Description=AWS EFS Mount Script

[Service]
ExecStart=/home/ubuntu/bin/aws-automount.sh

[Install]
WantedBy=multi-user.target
